# deploy/gpu-stt/systemd/lfg-stt-watchdog.service
[Unit]
Description=lfg STT watchdog (re-sync if box up but STT down)
[Service]
Type=oneshot
ExecStart=/home/dev/repos/lfg/deploy/gpu-stt/lifecycle.sh watchdog
