# deploy/gpu-stt/systemd/lfg-stt-stop.service
[Unit]
Description=Stop lfg STT GPU box (23:00 HKT)
[Service]
Type=oneshot
ExecStart=/home/dev/repos/lfg/deploy/gpu-stt/lifecycle.sh stop
