# deploy/parakeet-stt/parakeet-stt.service
[Unit]
Description=Parakeet STT host for lfg
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Adjust to wherever you checked this out / built the venv.
WorkingDirectory=/opt/parakeet-stt
EnvironmentFile=/opt/parakeet-stt/parakeet-stt.env
ExecStart=/opt/parakeet-stt/venv/bin/python /opt/parakeet-stt/server.py
Restart=on-failure
RestartSec=3
# First boot downloads the model from HF; give it room.
TimeoutStartSec=600

[Install]
WantedBy=multi-user.target
