# deploy/ops/lfg-tts-failover.timer
[Unit]
Description=Run lfg orb TTS failover each morning after the GPU box starts (03:10 UTC / 11:10 HKT)

[Timer]
# Box starts at 03:00 UTC (lfg-stt-start.timer); give it 10 min to boot + lazy-load.
OnCalendar=*-*-* 03:10:00
Persistent=true

[Install]
WantedBy=timers.target
