# deploy/ops/lfg-box-ensure.timer
[Unit]
Description=Retry lfg GPU box start every ~2 min until capacity frees, then fix orb TTS (self-disables on success)

[Timer]
# Fire shortly after boot (survives dev-box restart) and then every ~2 min.
OnBootSec=30
OnUnitActiveSec=120
Persistent=true

[Install]
WantedBy=timers.target
