# deploy/gpu-stt/lfg-stt-tunnel.service
[Unit]
Description=lfg STT SSH tunnel -> Novita 5090 (Cantonese SenseVoice on :8087)
After=network-online.target
Wants=network-online.target

[Service]
# Forward dev-box 127.0.0.1:8087 to the GPU box's 127.0.0.1:8087 over Novita's
# managed SSH. Germany (eu-ger-2) uses Novita's managed SSH which is PASSWORD
# auth only (it ignores container authorized_keys), so we use sshpass -e with
# SSHPASS from a 600 env file. systemd restarts if the channel drops.
# NOTE: host:port + password change if the instance is recreated/started —
# update ExecStart + ~/.config/lfg-stt-tunnel.env on redeploy.
EnvironmentFile=%h/.config/lfg-stt-tunnel.env
ExecStart=/usr/bin/sshpass -e /usr/bin/ssh -N -T \
  -o ServerAliveInterval=15 -o ServerAliveCountMax=3 \
  -o ExitOnForwardFailure=yes -o StrictHostKeyChecking=accept-new \
  -o UserKnownHostsFile=/home/dev/.ssh/known_hosts \
  -p 52248 \
  -L 127.0.0.1:8087:127.0.0.1:8087 \
  -L 127.0.0.1:8088:127.0.0.1:8088 \
  root@proxy.eu-ger-2.gpu-instance.novita.ai
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
